# school-bell.service
[Unit]
Description=校园智能广播定时系统
After=network.target sound.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/school-bell
ExecStart=/usr/bin/python3 /opt/school-bell/run.py
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:/opt/school-bell/logs/systemd.log
StandardError=append:/opt/school-bell/logs/systemd.log

# 安全加固
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/school-bell/data /opt/school-bell/logs /opt/school-bell/backups

[Install]
WantedBy=multi-user.target
